# headless/thunderbird-mcp.service
[Unit]
Description=Thunderbird (headless on Xvfb :99)
After=network-online.target

[Service]
Type=simple
Environment=DISPLAY=:99
ExecStart=/opt/thunderbird/thunderbird
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
